# own_files/files_9344_hi/etc/init.d/ser2net
#!/bin/sh /etc/rc.common
# Copyright (C) 2006-2011 OpenWrt.org

START=70

start () {
	stty -F /dev/ttyS0 57600
	ser2net -c /etc/ser2net.conf
}

stop() {
	killall ser2net
}

